# configs/hyprland/etc/profile.d/xcursors.sh
export XCURSOR_PATH=${XCURSOR_PATH}:~/.local/share/icons
